# ajapopaja-build.service
[Unit]
Description=Ajapopaja Build Service
After=docker.service
Requires=docker.service

[Service]
Restart=always
EnvironmentFile=/etc/ajapopaja-build.env
ExecStartPre=-/usr/bin/docker stop ajapopaja-build
ExecStartPre=-/usr/bin/docker rm ajapopaja-build
# Load the image in case it was updated or missing
ExecStartPre=/usr/bin/docker load -i /data/ajapopaja/docker/ajapopaja-build.tar
ExecStart=/usr/bin/docker run --name ajapopaja-build \
  -p ${PORT}:${PORT} \
  -v /home:/home \
  -v /Users:/Users \
  --add-host=host.docker.internal:host-gateway \
  -e MONGODB_URI=${MONGODB_URI} \
  -e OLLAMA_HOST=${OLLAMA_HOST} \
  -e DATABASE_NAME=${DATABASE_NAME} \
  -e PORT=${PORT} \
  -e WORKSPACES_ROOT=${WORKSPACES_ROOT} \
  ajapopaja-build:latest
ExecStop=/usr/bin/docker stop ajapopaja-build

[Install]
WantedBy=multi-user.target
